# systemd/hub-daily@.timer
[Unit]
Description=Run %i

[Timer]
Unit=hub-script@%i.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
