# systemd/hub-backup-appdata.service
[Unit]
Description=Backup Hub app data.
Requires=hub-services.service
After=hub-services.service

[Service]
Type=oneshot
ExecStart=/sbin/hub backup appdata
ExecStart=/sbin/hub notify "Hub updates" "Appdata backup successful!"
RemainAfterExit=false

[Install]
WantedBy=multi-user.target
