# systemd/hub-certbot.service
[Unit]
Description=Certbot to create/renew certificates from letsencrypt
Requires=hub-services.service
After=hub-services.service

[Service]
Type=oneshot
ExecStart=/usr/bin/hub cmd certbot
RemainAfterExit=false
ExecStartPost=/usr/bin/hub notify "Hub updates" "Certbot renewal successful!"

[Install]
WantedBy=multi-user.target
