# systemd/hub-certbot.timer
[Unit]
Description=Renew certbot certificates

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
