# systemd/hub-backup-appdata.timer
[Unit]
Description=Timer to back up hub app data daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
